# ansible/roles/voctomix/templates/systemd/videoteam-ingest.path.j2
# Ansible Managed

[Path]
PathExists=/dev/{{ item.1.devname }}

[Install]
WantedBy=paths.target
