# ansible/roles/voctomix/templates/systemd/videoteam-ingest.path.j2
# Ansible Managed

[Path]
PathExists={{ item.1.path_exists }}

[Install]
WantedBy=paths.target
